# crontab
0 0 * * * root python3 /app/pss-insight-data-sync.py
